# davmail/davmail.conf
description "davmail docker container"

start on started networking and started docker
stop on stopping networking or stopping docker

console none
respawn

exec docker start -a davmail

